# deployment/roles/sliver/templates/sliver.service.j2
[Unit]
Description=Sliver C2 Server

[Service]
ExecStart=/opt/sliver/sliver daemon

[Install]
WantedBy=multi-user.target